# supervisord.conf
[program:tornado]
command=/home/dotcloud/env/bin/python /home/dotcloud/current/main.py
